# res/osx/plist/com.cypherpunk.vpn.client.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.cypherpunk.vpn.client</string>
	<key>Program</key>
	<string>/Applications/CypherpunkVPN.app/Contents/MacOS/Electron</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Applications/CypherpunkVPN.app/Contents/MacOS/Electron</string>
	</array>
  <key>RunAtLoad</key>
	<true/>
</dict>
</plist>
